# deployment/ubuntu/wwjcar-backup.timer
[Unit]
Description=Daily WWJ Car Rent PostgreSQL backup

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
